# scripts/conduit.service
[Unit]
Description=Conduit daemon

[Service]
Type=forking
ExecStart=/usr/bin/conduit

[Install]
WantedBy=multi-user.target
